# sendmail-milter/supervisord.conf
[supervisord]
nodaemon = true

[program:dovecot]
process_name = dovecot
command = /usr/sbin/dovecot -F
startsecs = 0
autorestart = false

[program:sendmail]
process_name = sendmail-mta
directory = /etc/mail
command = sendmail-mta -Am -L sm-mta -bD -q10m -O LogLevel=20
startsecs = 0
autorestart = false

[program:rsyslog]
command=/usr/sbin/rsyslogd -n 
startsecs = 0
autorestart = false
